# fake.conf
[program: fake]
directory=/home/fake
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/fake/logfile_err.log
stdout_logfile=/home/fake/logfile.log
